# templates/etc/systemd/system/docker.service.d/service-overrides-archlinux.conf.erb
[Service]
EnvironmentFile=-/etc/conf.d/docker
