# templates/etc/systemd/system/docker.service.d/service-overrides-debian.conf.erb
[Service]
EnvironmentFile=-/etc/default/docker
EnvironmentFile=-/etc/default/docker-storage
ExecStart=
<%- if @docker_ce -%>
ExecStart=/usr/bin/<%= @docker_ce_command %> $OPTIONS \
<%- else -%>
ExecStart=/usr/bin/<%= @docker_command %> <%= @daemon_subcommand %> $OPTIONS \
<%- end -%>
        $DOCKER_STORAGE_OPTIONS
