# templates/etc/systemd/system/docker.service.d/service-overrides-debian.conf.erb
[Service]
EnvironmentFile=-/etc/default/docker
EnvironmentFile=-/etc/default/docker-storage
ExecStart=
ExecStart=/usr/bin/<%= @docker_daemon_command %> $OPTIONS \
        $DOCKER_STORAGE_OPTIONS
