# templates/etc/systemd/system/docker.service.d/service-overrides-debian.conf.erb
[Unit]
<% if @service_unit_after %>After=<%= @service_unit_after %><%end%>

[Service]
EnvironmentFile=-/etc/default/docker
EnvironmentFile=-/etc/default/docker-storage
ExecStart=
ExecStart=/usr/bin/<%= @docker_command %> <%= @daemon_subcommand %> $OPTIONS \
        $DOCKER_STORAGE_OPTIONS
