# templates/etc/systemd/system/docker.service.d/service-overrides-archlinux.conf.erb
[Unit]
<% if @service_unit_after %>After=<%= @service_unit_after %>

[Service]
EnvironmentFile=-/etc/conf.d/docker
