# templates/etc/systemd/system/docker.service.d/service-overrides-debian.conf.erb
[Service]
EnvironmentFile=-/etc/default/docker
EnvironmentFile=-/etc/default/docker-storage
<%- if @service_limits -%>
  <%- @service_limits.each do |lim,value| -%>
    <%- case lim -%>
      <%- when "nofile" -%>
LimitNOFILE=<%= value %>
      <%- when "nproc" -%>
LimitNPROC=<%= value %>
    <%- end -%>
  <%- end -%>
<%- end -%>
ExecStart=
ExecStart=/usr/bin/<%= @docker_command %> <%= @daemon_subcommand %> $OPTIONS \
        $DOCKER_STORAGE_OPTIONS
