# templates/etc/systemd/system/docker.service.d/service-overrides-debian.conf.erb
[Service]
EnvironmentFile=-/etc/default/docker
EnvironmentFile=-/etc/default/docker-storage
ExecStart=
<% unless @dockerd_daemon -%>
ExecStart=/usr/bin/<%= @docker_command %> <%= @daemon_subcommand %> $OPTIONS \
        $DOCKER_STORAGE_OPTIONS
<% end -%>
ExecStart=/usr/bin/<%= @dockerd_daemon %> $OPTIONS \
        $DOCKER_STORAGE_OPTIONS
