# templates/etc/init/docker-run.conf.erb
description "start and stop <%= title %> in docker"
author "Gareth Rushgrove"

start on (started docker)
stop on runlevel [!2345]

expect daemon

setuid root

respawn
respawn limit 5 20

exec docker run -d -u '<%= username %>' -h '<%= hostname %>'<% if dns %><% dns_array.each do |address| %> -dns <%= address %><% end %><% end %><% if env %><% env_array.each do |env| %> -e <%= env %><% end %><% end %><% if ports %><% ports_array.each do |port| %> -p <%= port %><% end %><% end %><% if volumes %><% volumes_array.each do |volume| %> -v <%= volume %><% end %><% end %><% if volumes_from %> -volumes-from <%= volumes_from %><% end %> -m <%= memory_limit %> <% if links %><%= links_array.map {|link| "-link #{link}" }.join(' ') %> <% end %><% if use_name %>-name <%= name %><% end %> <%= image %> <% if command %><%= command %><% end %>
