# templates/etc/init/docker-run.conf.erb
description "start and stop <%= title %> in docker"
author "Gareth Rushgrove"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root

respawn
respawn limit 5 20

exec docker run <% if port %>-p <%= port %> <% end %>-m <%= memory_limit %> <%= image %> <%= command %>
