# templates/etc/init/docker-run.conf.erb
description "start and stop <%= title %> in docker"
author "Gareth Rushgrove"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root

expect daemon

respawn
respawn limit 5 20

exec docker run -d -m <%= memory_limit %> <%= image %> <%= command %>
