# templates/etc/init/docker-run.conf.erb
# This file is managed by Puppet and local changes
# may be overwritten

description "start and stop <%= @title %> in docker"
author "Gareth Rushgrove"

start on (started <%= @service_name %>)
stop on stopping <%= @service_name %>

respawn
respawn limit 5 20

pre-start script
  cidfile=/var/run/docker-<%= @sanitised_title %>.cid
  if [ -f $cidfile ]; then
    cid="$(cat $cidfile)"
    if [ -n "$cid" ]; then
        if [ <%= @docker_command %> ps --no-trunc=true|grep -q $cid ]; then
            echo "Container $cid is still running.\n" >&2
            exit 1
        else
            rm -f $cidfile
        fi
    fi
fi
end script

script
<% if @use_name -%>
  <%= @docker_command %> inspect <%= @name %> && <%= @docker_command %> start -a -i <%= @name %>
  <%= @docker_command %> inspect <%= @name %> || <%= @docker_command %> run --cidfile=/var/run/docker-<%= @sanitised_title %>.cid \
<% else -%>
  <%= @docker_command %> run --cidfile=/var/run/docker-<%= @sanitised_title %>.cid \
<% end -%>
    <%= @docker_run_flags %> \
    <% if @use_name %>--name <%= @name %><% end %> \
    <% if @extra_parameters %><% @extra_parameters_array.each do |param| %>  \
      <%= param %> <% end %> \
    <% end -%> \
    <%= @image %> \
    <% if @command %> <%= @command %><% end %> &&
  exec <%= @docker_command %> wait "$(cat /var/run/docker-<%= @sanitised_title %>.cid)"
end script

post-stop script
if [ -e "/var/run/docker-<%= @sanitised_title %>.cid" ]; then
  <%= @docker_command %> kill "$(cat /var/run/docker-<%= @sanitised_title %>.cid)" && \
    rm "/var/run/docker-<%= @sanitised_title %>.cid"
fi
end script
