# templates/etc/init/docker-run.conf.erb
description "start and stop <%= @name %> in docker"
author "Gareth Rushgrove"

start on (started docker)
stop on stopping docker

setuid root

respawn
respawn limit 5 20

<% 
  env_hash = @env.is_a?(Hash) ? @env : Hash[@env.map{|kv| kv.split('=', 2)}]
  env_hash.each do |key,value|
-%>
env <%= key %>='<%= value.gsub(/'/, "'\"'\"'") %>'
<% end -%>

<%
  args = []

  args.concat ['-u', @username] if @username

  args.concat ['-h', @hostname] if @hostname

  args.concat ['-n', 'false'] if @disable_network

  @dns_array.each do |addr|
    args.concat ['-dns', addr]
  end

  env_hash.keys.each do |env_name|
    args.concat ['-e', env_name]
  end

  @ports_array.each do |port|
    args.concat ['-p', port]
  end

  @volumes_array.each do |volume|
    args.concat ['-v', volume]
  end

  @volumes_from.each do |volume_from|
    args.concat ['--volumes-from', volume_from]
  end

  args.concat ['-m', @memory_limit]

  @links_array.each do |link|
    args.concat ['--link', link]
  end

  args.concat ['-name', @name] if @use_name

  @lxc_conf_array.each do |lxc_conf|
    args.concat ['--lxc-conf', lxc_conf]
  end

  args << @image

  args << @command if @command

  args.map!(&:to_s) # fix a ruby 1.9.3 issue where shellwords can't handle symbols
  require 'shellwords'
-%>

exec docker run <%= Shellwords.shelljoin(args) %>
